# systemd/dot-config/systemd/user/wm-session.target
[Unit]
Description=mangowm Session Target
Requires=graphical-session.target
After=graphical-session.target
